# recipes-domu/domu-image-fusion/files/meta-xt-prod-extra/recipes-aos/aos-servicemanager/files/aos.target
[Unit]
Description=Aos services
Wants=aos-servicemanager.service aos-updatemanager.service
After=aos-servicemanager.service aos-updatemanager.service

[Install]
WantedBy=multi-user.target
