# recipes-domu/domu-image-fusion/files/meta-xt-prod-extra/recipes-aos/aos-servicemanager/files/first-boot.service
[Unit]
Description=Runs first_boot.sh once at first boot

[Service]
Type=simple
ExecStart=/usr/bin/first_boot.sh

[Install]
WantedBy=multi-user.target
